# islom.conf
[program:islombot]
directory=/home/ubuntu/islomuzbot-main
command=python3 app.py

autostart=true
autorestart=true
environment=HOME="/home/ubuntu",USER="ubuntu"

stderr_logfile=/home/ubuntu/islomuzbot-main/logfile_err.log
stdout_logfile=/home/ubuntu/islomuzbot-main/logfile.log